# scripts/com.joelgerard.garmindb.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE PLIST PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.joelgerard.garmindb</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/joelgerard/dev/git/GarminDB/scripts/daily_sync.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardErrorPath</key>
    <string>/Users/joelgerard/HealthData/logs/sync_error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/joelgerard/HealthData/logs/sync_out.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/joelgerard/dev/git/GarminDB</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
